# Neo.Node.Service/neo-node.service
[Unit]
Description=NEO N3 blockchain node

[Service]
Type=notify
ExecStart=/usr/bin/neo-node

[Install]
WantedBy=multi-user.target
